# recipes-remidomo/overrides/udev/50-update.rules
ACTION=="add", KERNEL=="sd[a-z][0-9]", RUN+="/usr/bin/update.sh"
